# etc/system/lunch-scraper.timer
[Unit]
Description=Run Lunch Menu Scraper hourly
Requires=lunch-scraper.service

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=300
AccuracySec=1min

[Install]
WantedBy=timers.target
